# packaging/deb-vault/deb-pkg/lib/systemd/system/vault.service
[Unit]
Description=Vault
Requires=consul-online.target
After=consul-online.target

[Service]
Restart=on-failure
ExecStart=/usr/local/bin/vault server -config /etc/vault
ExecReload=/bin/kill -HUP $MAINPID
KillSignal=SIGTERM
User=vault
Group=vault

[Install]
WantedBy=multi-user.target
