# packaging/deb-consul-bootstrap-aws/deb-pkg/etc/systemd/system/consul.service.d/consul.conf
[Service]
ExecStartPre=/usr/local/bin/consulBootstrap.sh
PermissionsStartOnly=true
